# eos-chrome-plugin-update.service.in
[Unit]
Description=Update or install Google Chrome plugins

[Service]
Type=simple
ExecStart=@pkglibexecdir@/eos-chrome-plugin-update
