# eos-pepflash-update.service.in
[Unit]
Description=Update or install Pepper Flash plugin

[Service]
Type=simple
ExecStart=@pkglibexecdir@/eos-pepflash-update
